# bin/com.pavpan.top-cpu-widget.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pavpan.top-cpu-widget</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/pavpan/bin/top-cpu-widget</string>
    <string>--window-minutes</string>
    <string>60</string>
    <string>--sample-seconds</string>
    <string>1</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/top-cpu-widget.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/top-cpu-widget.err.log</string>
</dict>
</plist>
